# postgres.container
# postgres.container - Podman Quadlet for PostgreSQL
# Save this file to ~/.config/containers/systemd/postgres.container

[Container]
# PostgreSQL image to use
Image=docker.io/library/postgres:17

# Environment variables for PostgreSQL setup
Environment=POSTGRES_USER=postgres
Environment=POSTGRES_PASSWORD=postgres
Environment=POSTGRES_DB=postgres
Environment=PGDATA=/var/lib/postgresql/data/pgdata

# Expose PostgreSQL port
PublishPort=5432:5432

# Volume mappings for persistent data
Volume=postgres-data:/var/lib/postgresql/data

# Set resource limits
CPUQuota=100%
MemoryLimit=1G

# Restart policy
Restart=always

# Health check
HealthCmd=["pg_isready", "-U", "postgres"]
HealthInterval=10s
HealthTimeout=5s
HealthRetries=3
HealthStartPeriod=30s

# Labels
Label=description="PostgreSQL Database Server"
Label=app=postgres

[Install]
# Start at boot
WantedBy=default.target
